# test/wstest.conf
[program:wstest]

command     = /srv/awesome/test/webtest.py
directory   = /srv/awesome/test
user        = root
startsecs   = 3

redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
